# debian/rails-new.cron.d
* * * * * rails-new cd /usr/lib/rails-new; /usr/bin/rails-production foo
